# src/deb/control/prerm
#!/bin/bash

service indexmic stop
if [ -d /run/systemd/system ];
then
  systemctl disable indexmic.service
fi
